# final/sh_timer.service
[Unit]
Description=My script

[Service]
ExecStart='/home/qym/Information-Retrieval/final/crawler.sh'

[Install]
WantedBy=multi-user.target
