# final/sh_timer.timer
[Unit]
Description=Run crawler every day at midnight

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
